# mongo/crontab
* * * * * echo "Hello world" >> /dev/stdout 2>&1
# crontab requires empty line at end of file